# bot_sites/systemd/mirrorhub-bot.service
[Unit]
Description=MirrorHub Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/root/mirrorhub
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/mirrorhub/.venv/bin/python -m bot.main
Restart=always

[Install]
WantedBy=multi-user.target
